# profile.d/home-local-bin-path.sh
#!/bin/sh
# add path to user's local bin
append_path "$HOME/.local/bin"
export PATH
